# src/cli/service/templates/proton-drive-sync.service
[Unit]
Description=Proton Drive Sync
After=network-online.target proton-drive-sync-keyring.service
Wants=network-online.target
Requires=proton-drive-sync-keyring.service

[Service]
Type=simple
{{USER_LINE}}
ExecStart={{BIN_PATH}} start --no-daemon
Restart=on-failure
RestartSec=5
Environment="HOME={{HOME}}"
Environment="PATH={{HOME}}/.local/bin:/usr/local/bin:/usr/bin:/bin"
Environment="DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/{{UID}}/bus"
EnvironmentFile=-{{KEYRING_ENV_FILE}}

[Install]
WantedBy={{WANTED_BY}}
